# docker/containers/application/config/system/rc.local
#!/bin/bash

service php7.3-fpm start
# service varnish start
service nginx start
# service cron start
